# .crontab
# 1. minute (0-59)
# |   2. hour (0-23)
# |   |   3. day of month (1-31)
# |   |   |   4. month (1-12)
# |   |   |   |   5. day of week (0-7: 0 or 7 is Sun, or use names)
# |   |   |   |   |   6. commandline
# |   |   |   |   |   |
#min hr dom mon dow command

@reboot rm -rf "$HOME/.cache/" 2>/dev/null
@reboot rm -rf "$HOME/.thumbnails" 2>/dev/null
@reboot rm -rf "$HOME/.local/share/Trash" 2>/dev/null
@reboot rm -rf $HOME/.wine/drive_c/windows/temp/* 2>/dev/null
